# docker-files/php-queue-worker/conf.d/acme-worker.conf
[program:acme-queue]
process_name = %(program_name)s_%(process_num)02d
command = php artisan queue:work queue --delay=0 --memory=256 --timeout=300 --sleep=1 --tries=1
directory = /var/www/acme
user = captain
numprocs = 2
stdout_logfile = /var/www/acme/storage/logs/queue.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 20
redirect_stderr = true
autostart = true
autorestart = true
